# etc/scripts/running/eard.service.in
[Unit]
Description=EARD - Energy Aware Runtime node daemon
ConditionPathExists=@tmpdir@

[Service]
Type=simple
EnvironmentFile=@destdir@/etc/slurm/ear.conf
EnvironmentFile=@destdir@/etc/slurm/ear.link
ExecStart=@destdir@/sbin/eard 1 $EAR_TMP 0
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
RestartSec=30s
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target

