# etc/systemd/eargmd.service.in
[Unit]
Description=EARGMD - EAR Global Manager
ConditionPathExists=@localstatedir@
# Requires=mariadb.service
# Wants=slurmd.service
# After=mariadb.service,slurmd.service

[Service]
Type=simple
Environment="ETC=@sysconfdir@"
# ExecStart also accepts the verbosity level as first parameter
ExecStart=@prefix@/sbin/eargmd
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target

